# server/openclaw.service
# Template — copy to server: scp openclaw.service root@<serverIp>:/etc/systemd/system/
# Requires: /root/.openclaw/.env  (see server/.env.example)
# Enable:   systemctl daemon-reload && systemctl enable --now openclaw
# Logs:     journalctl -u openclaw -f -n 100
# Access:   https://<hostname>.<tailnet>/  (via Tailscale Serve)

[Unit]
Description=OpenClaw Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw

# Secrets + model tiers
EnvironmentFile=/root/.openclaw/.env

# Runtime environment
Environment=NODE_ENV=production
Environment=HOME=/root/.openclaw
Environment=XDG_CONFIG_HOME=/root/.openclaw
Environment=OPENCLAW_GATEWAY_PORT=18789
Environment=TERM=xterm-256color

WorkingDirectory=/root/.openclaw

# TODO: replace with actual binary from `docker inspect` (see plan Step 1.1)
ExecStart=/usr/bin/openclaw gateway

Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/root/.openclaw
ProtectHome=false
PrivateTmp=true

# Logging → journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw

[Install]
WantedBy=multi-user.target
